# air.service
[Unit]
Description=Air

[Service]
ExecStart=/bin/ruby '/home/pi/do_record_data.rb'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
